# 01-configure/roles/photoprism/templates/rclone.service.j2
[Unit]
Description=Download latest pictures using rclone

[Service]
Type=simple
ExecStart={{ ansible_user_dir }}/rclone.sh
